# docker/crontab
# 每週三 10:00（Asia/Taipei，由容器 TZ 決定）
0 10 * * 3 /bin/bash /app/scripts/weekly_sync.sh
