# scripts/com.crash-trend.weekly-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- launchd 備選排程（不用 Docker 時）：把兩處路徑改成你的安裝位置後，
     複製為 ~/Library/LaunchAgents/com.crash-trend.weekly-sync.plist 並 launchctl load -->
<dict>
  <key>Label</key>
  <string>com.crash-trend.weekly-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/PATH/TO/crash-trend/scripts/weekly_sync.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>37</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
